# deploy/osrs-sync-ingest.service
[Unit]
Description=OSRS Sync Ingest Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/<user>/osrs-sync-ingest
EnvironmentFile=/home/<user>/osrs-sync-ingest/.env
ExecStart=/home/<user>/osrs-sync-ingest/.venv/bin/uvicorn app:app --host ${HOST} --port ${PORT}
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
